# home/config/systemd/user/sync-user-systemd.service
[Unit]
Description=Sync systemd files in dotfiles

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sync-etc-systemd.sh '/home/aetf/customizations/config/systemd/user' '/home/aetf/systemd/user'
