# home/config/systemd/user/sync-user-systemd.path
[Unit]
Description=Path monitor for syncing systemd files in dotfiles

[Path]
PathChanged=/home/aetf/customizations/config/systemd/user

[Install]
WantedBy=multi-user.target
