# home/config/systemd/user/sync-user-systemd.path
[Unit]
Description=Path monitor for syncing systemd files in dotfiles

[Path]
PathChanged=/home/peifeng/customizations/home/config/systemd/user

[Install]
WantedBy=default.target
