# .config/aconfmgr/files/etc/systemd/system/mnt-Aetf\x2dLaptop.automount
[Unit]
Description=Automount /mnt/Aetf-Laptop

[Automount]
Where=/mnt/Aetf-Laptop

[Install]
WantedBy=multi-user.target
