# .config/aconfmgr/files/etc/udev/rules.d/92-scsi-sg.rules
# load SCSI generic (sg) driver for disk enclosure
SUBSYSTEM=="scsi", DRIVERS=="ses", RUN+="/sbin/modprobe -bv sg"
