# etc/profile.d/apache-spark.sh
#!/bin/sh
SPARK_HOME=/usr/share/apache-spark
PYSPARK_PYTHON=python3
export SPARK_HOME
