# etc/systemd/system/refind-upgrade.path
[Unit]
Description=Path monitor for rEFInd upgrades

[Path]
PathChanged=/usr/share/refind

[Install]
WantedBy=multi-user.target
