# etc/systemd/system/refind-upgrade.service
[Unit]
Description=rEFInd boot manager upgrade

[Service]
Type=oneshot
ExecStart=/usr/local/bin/refind-upgrade
